# podman/quadlets/n8n-main.container
[Unit]
Description=N8N main application
After=n8n-postgres.service n8n-redis.service

[Container]
Image=n8nio/n8n:latest
ContainerName=n8n-main
PublishPort=5678:5678
EnvironmentFile=REPO_DIR_PLACEHOLDER/.env
Environment=DB_TYPE=postgresdb
Environment=DB_POSTGRESDB_HOST=postgres
Environment=DB_POSTGRESDB_PORT=5432
Environment=DB_POSTGRESDB_DATABASE=n8n
Environment=DB_POSTGRESDB_USER=n8n
Environment=EXECUTIONS_MODE=queue
Environment=QUEUE_BULL_REDIS_HOST=redis
Environment=QUEUE_BULL_REDIS_PORT=6379
Environment=QUEUE_HEALTH_CHECK_ACTIVE=true
Environment=N8N_PORT=5678
Environment=N8N_PROTOCOL=http
Volume=n8n-data.volume:/home/node/.n8n:Z
Network=n8n-network.network
UserNS=keep-id:uid=1000,gid=1000
Label=app=n8n

[Service]
Restart=always
ExecStartPre=/bin/sh -c 'until podman exec postgres pg_isready -U n8n; do sleep 2; done'
ExecStartPre=/bin/sh -c 'until podman exec redis redis-cli ping; do sleep 2; done'

[Install]
WantedBy=default.target
